# assets/environment-interpreter/common/etc/profile.d/0502_jupyter.sh
# shellcheck shell=bash
# ============================================================================ #
#
# Setup Jupyter
#
# ---------------------------------------------------------------------------- #

_jupyter_env="${FLOX_ENV}/share/jupyter"
if [[ -d "$_jupyter_env" ]]; then
  export JUPYTER_PATH="${_jupyter_env}${JUPYTER_PATH:+:$JUPYTER_PATH}"
fi

# ---------------------------------------------------------------------------- #
#
#
#
# ============================================================================ #
