# assets/environment-interpreter/common/etc/profile.d/0503_java.sh
# shellcheck shell=bash
# ============================================================================ #
#
# Setup Java
#
# ---------------------------------------------------------------------------- #

if [[ -x "$FLOX_ENV/bin/java" ]]; then
  export JAVA_HOME="$FLOX_ENV"
fi
# ---------------------------------------------------------------------------- #
#
#
#
# ============================================================================ #
